# ui/macos/com.secsec.ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the secsec menu-bar app. Install to ~/Library/LaunchAgents/com.secsec.ui.plist,
  then: launchctl load -w ~/Library/LaunchAgents/com.secsec.ui.plist
  RunAtLoad starts it at login (which is when it pops the passphrase prompt). Adjust the path below
  if you installed the .app somewhere other than /Applications.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.secsec.ui</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/secsec-menubar.app/Contents/MacOS/secsec-menubar</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>ProcessType</key>
	<string>Interactive</string>
</dict>
</plist>
